# deploy/systemd/pacd-healthcheck.service
[Unit]
Description=Pingancoin pacd health check

[Service]
Type=oneshot
EnvironmentFile=-/etc/pingancoin/pacd-healthcheck.env
ExecStart=/usr/local/bin/pacd-health-check
